# deploy/scimtest-server.service
[Unit]
Description=scimtest application tunnel server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=true
EnvironmentFile=/etc/scimtest-server/scimtest-server.env
ExecStart=/usr/local/bin/scimtest-server --addr 127.0.0.1:8000 --domain scimtest.rselbach.com --dashboard-domain admin.scimtest.rselbach.com --scheme https --behind-proxy --data /var/lib/scimtest-server/scimtest-server.json --logs
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
StateDirectory=scimtest-server

[Install]
WantedBy=multi-user.target
